# docker/epics-gateway/copyroot/etc/profile.d/10-epics-gateway.sh
export PATH="/EPICS/extensions/bin/linux-x86_64:${PATH}"

export GATEWAY_HOME="/gateway"
export EPICS_CA_ADDR_LIST="172.17.255.255"
export EPICS_CAS_INTF_ADDR_LIST=""

